# .local/skel/etc/cron.weekly/btrfs-subv-backup
#!/bin/sh

do_btrfs-subv-backp weekly  4 


