# .local/skel/etc/cron.weekly/btrfs-snaps-home
#!/bin/sh

btrfs-snaps-home weekly  4 | egrep -v '(Create a snapshot of|Will delete the oldest|Delete subvolume|Making snapshot of )' || true


