# skel/etc/cron.weekly/btrfs-snaps
#!/bin/sh

btrfs-snaps weekly  4 | egrep -v '(Create a snapshot of|Will delete the oldest|Delete subvolume|Making snapshot of )' || true


